# deploy/fund-rate-scan.service
[Unit]
Description=Funding Rate Signal Scanner
After=network.target

[Service]
Type=simple
WorkingDirectory=<home>/fund_rate_arb
ExecStart=<home>/.local/bin/uv run python -m fund_rate_arb.main
Restart=on-failure
RestartSec=30
StandardOutput=append:<home>/fund_rate_arb/logs/scan.log
StandardError=append:<home>/fund_rate_arb/logs/scan.log
EnvironmentFile=<home>/fund_rate_arb/.env

[Install]
WantedBy=default.target
