# deploy/systemd.service
[Unit]
Description=Funding Rate Arbitrage Signal Scanner
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/fund_rate_arb
ExecStart=/opt/fund_rate_arb/.venv/bin/python -m fund_rate_arb.main
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
